# autostart/xrdb.desktop
[Desktop Entry]
Type=Application
Name=X resources load
Exec=xrdb -load /home/dacheng/.Xresources
Comment=Load X resources
