# supervisord.d/yii-queue-worker.ini
[program:yii-queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /path/to/yii queue/listen --color=0 --verbose=1
autostart=1
autorestart=1
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/path/to/yii-queue-worker.log
